# tools/specs/linux/systemd.yaml
# systemd unit health for Linux hosts scraped by node_exporter.
uid: linux-systemd
title: Linux / systemd
description: >-
  systemd unit health for Linux hosts scraped by node_exporter: the count of failed
  units, the active/inactive breakdown, and a live table of exactly which units have
  failed and where. Answers "is any service down on these hosts, and which one?"
  rather than leaving you to SSH and run systemctl.
datasource: prometheus
tags: [linux, systemd, services, node-exporter]
seo_phrase: Node Exporter systemd Grafana dashboard
refresh: 30s
time_from: now-6h

questions:
  - "Are any systemd units in the failed state right now, and on which hosts?"
  - "Which specific units failed (the table you would otherwise get from systemctl)?"
  - "Is the failed-unit count trending up across the fleet?"
  - "How many units are active vs inactive overall?"

production_lessons: >-
  A failed systemd unit is the most direct "something is broken" signal a host can emit, yet
  it is invisible on CPU/memory/disk dashboards — a crashed exporter, a dead log shipper or
  a failed mount unit leaves the box looking perfectly healthy. So this dashboard leads with
  a single red **failed-unit count** and backs it with a **table of the actual unit names**
  so on-call can see `nginx.service on web-03` without opening a shell. The subtle gotcha is
  that node_exporter only reports the units it can see under its allow/deny lists; if a unit
  you care about never appears, it is being filtered at the collector, not magically
  healthy. Treat any non-zero failed count as actionable and let the table tell you who.

data_sources:
  - "`node_exporter` `systemd` collector (`node_systemd_unit_state`, `node_systemd_units`).
    Enable it with `--collector.systemd`; it is off by default."

templating:
  - name: job
    label: Job
    query: label_values(node_uname_info, job)
    description: Prometheus scrape job for your node_exporter targets.
    multi: false
    includeAll: false
  - name: instance
    label: Instance
    query: label_values(node_uname_info{job="$job"}, instance)
    description: Host(s) to display; supports multi-select.
    multi: true
    includeAll: true

rows:
  - title: Service health — read these first
    panels:
      - type: stat
        title: Failed units
        description: Total systemd units in the failed state across selected hosts. Any value above zero is actionable.
        unit: short
        decimals: 0
        colorMode: background
        thresholds:
          - {color: green, value: null}
          - {color: red, value: 1}
        targets:
          - expr: sum(node_systemd_unit_state{job="$job", instance=~"$instance", state="failed"} == 1) or vector(0)
            instant: true
            legend: failed
      - type: stat
        title: Hosts with failures
        description: Number of distinct hosts that have at least one failed unit.
        unit: short
        decimals: 0
        colorMode: background
        thresholds:
          - {color: green, value: null}
          - {color: red, value: 1}
        targets:
          - expr: >-
              count(count by (instance) (node_systemd_unit_state{job="$job", instance=~"$instance", state="failed"} == 1)) or vector(0)
            instant: true
            legend: hosts
      - type: stat
        title: Active units
        description: Total units currently in the active state across selected hosts.
        unit: short
        decimals: 0
        colorMode: value
        thresholds:
          - {color: green, value: null}
        targets:
          - expr: sum(node_systemd_unit_state{job="$job", instance=~"$instance", state="active"} == 1) or vector(0)
            instant: true
            legend: active
      - type: stat
        title: Inactive units
        description: Total units in the inactive state — usually fine (stopped one-shots and sockets), shown for context.
        unit: short
        decimals: 0
        colorMode: value
        thresholds:
          - {color: green, value: null}
        targets:
          - expr: sum(node_systemd_unit_state{job="$job", instance=~"$instance", state="inactive"} == 1) or vector(0)
            instant: true
            legend: inactive

  - title: Which units failed
    panels:
      - type: table
        title: Failed units
        description: Live list of every unit in the failed state with its host. This is the systemctl --failed view, fleet-wide.
        unit: short
        height: 10
        targets:
          - expr: node_systemd_unit_state{job="$job", instance=~"$instance", state="failed"} == 1
            instant: true
            format: table
            legend: "{{instance}} {{name}}"

  - title: Trend
    panels:
      - type: timeseries
        title: Failed units over time (per host)
        description: Per-host failed-unit count. A step up marks the moment a service died; a step down marks recovery.
        unit: short
        min: 0
        drawStyle: bars
        legend: table
        targets:
          - expr: sum by (instance) (node_systemd_unit_state{job="$job", instance=~"$instance", state="failed"} == 1)
            legend: "{{instance}}"
      - type: timeseries
        title: Active vs inactive vs failed (fleet)
        description: Fleet-wide unit-state breakdown over time for context on the failed count.
        unit: short
        min: 0
        legend: table
        targets:
          - expr: sum(node_systemd_unit_state{job="$job", instance=~"$instance", state="active"} == 1)
            legend: active
          - expr: sum(node_systemd_unit_state{job="$job", instance=~"$instance", state="inactive"} == 1)
            legend: inactive
          - expr: sum(node_systemd_unit_state{job="$job", instance=~"$instance", state="failed"} == 1)
            legend: failed

alerts:
  - name: HostSystemdUnitFailed
    severity: critical
    expr: node_systemd_unit_state{state="failed"} == 1
    for: 5m
    summary: "Unit {{ $labels.name }} failed on {{ $labels.instance }}"
    why: A failed unit means a service the host is supposed to run is down — and it stays down until something restarts it.
    investigate: "On the host run `systemctl status {{ $labels.name }}` and `journalctl -u {{ $labels.name }} -e` to see the failure reason."
    runbook: "Fix the root cause (config, dependency, resource), then `systemctl reset-failed {{ $labels.name }} && systemctl restart {{ $labels.name }}`."
    recovery: Clears when the unit leaves the failed state (restarted or reset).
    false_positives: One-shot units that fail by design in some environments — exclude specific unit names or use `state="failed"` only for service units you own.
  - name: HostManySystemdFailures
    severity: warning
    expr: sum by (instance, job) (node_systemd_unit_state{state="failed"} == 1) > 3
    for: 5m
    summary: "Host {{ $labels.instance }} has more than 3 failed units"
    why: Several units failing together usually points at a shared root cause — a full disk, a missing mount, a failed dependency, or a bad config rollout.
    investigate: "Open the failed-units table for the host and look for a common dependency (a mount, network target, or shared config)."
    runbook: "Fix the shared cause first (e.g. remount the failed filesystem) then reset-failed the dependent units; they often recover together."
    recovery: Clears when the host's failed-unit count drops to 3 or fewer for 5m.
    false_positives: A host with several intentionally-disabled one-shots — tune the threshold or exclude those units.

troubleshooting:
  - symptom: All systemd panels show "No data"
    cause: The systemd collector is not enabled (it is off by default).
    action: Start node_exporter with `--collector.systemd`; confirm `node_systemd_unit_state` appears in Explore.
  - symptom: A unit you care about never shows up
    cause: The collector's allow/deny lists (`--collector.systemd.unit-include/exclude`) filter it out.
    action: Adjust the include regex so the unit is exported; absence is not the same as healthy.
  - symptom: Failed count flaps between 0 and 1
    cause: A unit in a crash-restart loop oscillating through failed and activating.
    action: Inspect `systemctl status` for restart counts; fix the crash rather than the symptom.

performance: >-
  The systemd collector can export thousands of series per host if no include/exclude filter
  is set, which is the main cost driver here. Scope it at the collector with
  `--collector.systemd.unit-include` to the units you actually run. All panels use cheap
  `sum`/`count` over the boolean `== 1` states, so query cost scales with that series count.

customization: >-
  Restrict the exported units at the collector to cut cardinality, then this dashboard
  needs no PromQL changes. Lower the "many failures" threshold for small hosts, or split the
  table by unit `type` (service/mount/timer) if you want separate views for each.
